# planetlab/ipfw.cron
# Runs every 5 minutes and clean ipfw expired rules
# $Id: ipfw.cron 6069 2010-04-15 09:35:33Z marta $
*/5 * * * * root     echo "super killexpired" | /vsys/ipfw-be root > /dev/null 2>&1
